# trustserver/trunk/etc/gentoo/konfidi-trustserver
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting konfidi trustserver"
        start-stop-daemon --start --quiet --background --chuid konfidi:konfidi --chdir /opt/konfidi/trustserver --exec /opt/konfidi/trustserver/TrustServer.py
	eend $?
}

stop() {
	ebegin "Stopping konfidi trustserver"
	start-stop-daemon --stop --quiet --pidfile /opt/konfidi/trustserver/run/trustserver.pid
	eend $?
}
